# files/systemd/user/mako.service.d/10-niri.conf
[Unit]
PartOf=niri.service

[Install]
WantedBy=niri.service
